# root/etc/services.d/psx-docker/run
#!/usr/bin/with-contenv bash
python3 /opt/psx-docker/start.py